# contrib/dev.modelrouter.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>dev.modelrouter</string>

    <key>ProgramArguments</key>
    <array>
        <string>/usr/local/bin/modelrouter</string>
        <string>serve</string>
    </array>

    <key>EnvironmentVariables</key>
    <dict>
        <key>MODELROUTER_CONFIG</key>
        <string>${HOME}/.modelrouter/config.toml</string>
    </dict>

    <key>RunAtLoad</key>
    <true/>

    <key>StandardOutPath</key>
    <string>${HOME}/.modelrouter/modelrouter.log</string>

    <key>StandardErrorPath</key>
    <string>${HOME}/.modelrouter/modelrouter.log</string>

    <key>KeepAlive</key>
    <true/>
</dict>
</plist>
